# supervisor.conf
[program:rundeck]
command=java -XX:MaxPermSize=256m -Xmx1024m -Dserver.http.port=80 -Dserver.https.port=443 -jar /rundeck.jar
redirect_stderr=true
